# provision/config_files/system-update.service
[Unit]
Description=Full system update & reboot

[Service]
Type=oneshot
ExecStart=/bin/pacman -Syu --noconfirm
ExecStartPost=/bin/reboot
